# calendar.conf
start on started cmd
stop on stopping cmd

export LANG LC_ALL

chdir /mnt/us/extensions/calendar 
exec ./update.sh 

